# roms/skiboot/external/opal-prd/opal-prd.service
[Unit]
Description=OPAL PRD daemon
ConditionVirtualization=false
ConditionPathExists=/sys/firmware/devicetree/base/ibm,opal/diagnostics

[Service]
ExecStart=/usr/sbin/opal-prd
Restart=always

[Install]
WantedBy=multi-user.target
